# QA/etc/crons/.crontab
*/5 * * * * /Users/egorbedov/Downloads/Prog/TechnoPark/Web/Q\&A/TechnoPark-Web/QA/etc/crons/top_users.sh